# playbooks/roles/sv_user_env/templates/profile.d/proxy.csh.j2
# proxy env
{% for item in proxy_env | dict2items %}
setenv {{ item.key }} "{{ item.value }}"
{% endfor %}
